# crontab
0 8,14 * * * /home/colejhudson/.colejhudson/dotfiles/scripts/cron/post-ip.sh
